# assets/madsonic.conf
[program:madsonic]
user = root
command = /opt/madsonic/start.sh
